# system/firstsetup/skel/.config/autostart/firstsetup.desktop
[Desktop Entry]
Type=Application

Name=firstsetup
Exec=/etc/firstsetup/bin/firstsetup
Icon=preferences-desktop

Comment=First Setup