# services/nix.mount
[Unit]
Description=Mount nix-store at /nix

[Mount]
What=/etc/nix-store
Where=/nix
Options=bind,x-mount.mkdir

[Install]
WantedBy=local-fs.target