# services/postinstall.service
[Unit]
Description=Run postinstall setup when network connection is detected
After=network.target

[Service]
Type=oneshot
ExecStart=/bin/bash -c 'while ! ping -c 1 8.8.8.8 >/dev/null 2>&1; do sleep 5; done; /usr/bin/postinstall'

[Install]
WantedBy=multi-user.target
