# services/systemsetup.service
[Unit]
Description=Setup system
After=network.target

[Service]
Type=oneshot
ExecStart=/usr/bin/blackbox-terminal -c '/usr/bin/bash /etc/firstsetup/bin/setup'
RemainAfterExit=true
Environment=USER=%I
ExecStartPost=/usr/bin/systemctl disable systemsetup.service

[Install]
WantedBy=multi-user.target
