# services/systemsetup.service
[Unit]
Description=System setup preperation
After=gdm.service

[Service]
Type=oneshot
ExecStart=/usr/bin/bash /etc/firstsetup/bin/prep
RemainAfterExit=true
ExecStartPost=/usr/bin/systemctl disable systemsetup.service

[Install]
WantedBy=default.target
