# services/systemsetup.service
[Unit]
Description=System setup preperation
After=gdm.service

[Service]
Type=oneshot
ExecStart=/usr/bin/mv /etc/firstsetup/misc/setup.desktop /home/*/.config/autostart
RemainAfterExit=true
ExecStartPost=/usr/bin/systemctl diable systemsetup.service

[Install]
WantedBy=default.target
